# ansible/roles/speedtesting/templates/speedtest.service.j2
[Service]
Type=oneshot
ExecStart=/usr/local/bin/pipenv run python speedtest.py --influx-database speedtesting
WorkingDirectory=/opt/speedtesting
User=speedtesting
